# systemd/user/mopidy.service
[Unit]
Description=Music Player Daemon

[Service]
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=default.target
